# ansible/configure-tunneled-raspberry-pi/vanilla-klipper-raspberry-pi-5/files/ports.service
[Unit]
Description=USB Serial Bridge Script
After=network.target syslog.target

[Service]
Type=simple
ExecStart=/opt/ports.sh
User=root
WorkingDirectory=/opt
Restart=on-failure

[Install]
WantedBy=multi-user.target